# files/system/usr/lib/systemd/system/onboot.service
[Unit]
Description=onboot

[Service]
Type=oneshot
ExecStart=/usr/bin/onboot.sh

[Install]
WantedBy=default.target
