# files/system/usr/lib/systemd/system/var-mnt-sscd.mount
[Unit]
Description=Mount sscd

[Mount]
What=/dev/sr0
Where=/var/mnt/sscd
Type=iso9660
Options=ro,noatime

[Install]
WantedBy=multi-user.target
